# crontab
#Note: Depending on the OS the path may need to be full 
* * * * * ~/bitcoin_env/bin/python3 ~/bitcoin_code/harvester.py
